# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=www-data

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off; pid /tmp/nginx.pid;"
autostart=true
autorestart=true
priority=20
stopsignal=QUIT
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ==========================================================================
# Email queue processor — processes pending emails every 60 seconds
# Uses the Symfony Console command (bin/console email:process-queue)
# ==========================================================================
[program:email-queue]
command=/bin/sh -c "BACKOFF=60; while true; do if php /var/www/bin/console email:process-queue --batch-size=50 --reminders 2>&1; then BACKOFF=60; else BACKOFF=$((BACKOFF * 2)); [ \"$BACKOFF\" -gt 900 ] && BACKOFF=900; echo \"[email-queue] echec, retry dans ${BACKOFF}s\"; fi; sleep $BACKOFF; done"
autostart=true
autorestart=true
priority=30
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ==========================================================================
# Monitoring — collects metrics, evaluates thresholds, sends alerts every 5 min
# ==========================================================================
[program:monitoring]
command=/bin/sh -c "sleep 30; while true; do php /var/www/bin/console monitor:check --cleanup 2>&1; sleep 300; done"
autostart=true
autorestart=true
priority=35
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ==========================================================================
# Rate limit cleanup — purges expired rate-limit files every 30 minutes
# ==========================================================================
[program:ratelimit-cleanup]
command=/bin/sh -c "while true; do php /var/www/bin/console ratelimit:cleanup 2>&1; sleep 1800; done"
autostart=true
autorestart=true
priority=40
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
